# pkgdists/slurm/base/etc/supervisord.d/munged.conf
[program:munged]
command=/usr/sbin/munged -F
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
