# pkgdists/slurm/slurmd/etc/supervisord.d/slurmd.conf
[program:slurmd]
command=/usr/local/sbin/slurmd -D -v
environment=LD_LIBRARY_PATH=/usr/local/lib/
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
#stdout_logfile=syslog
